# rootfs/etc/service/php/run
#!/usr/bin/env bash

php-fpm -y ${PHP_INI_DIR}/php-fpm.conf